# govps-cleanlog.timer
[Unit]
Description=govps cleanlog timer

[Timer]
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
