# assets/systemd/goproxy-cleanlog.timer
[Unit]
Description=goproxy cleanlog timer

[Timer]
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
